# core/src/main/resources/tomcat/service/solaris/tomcat9.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<!--  Service manifest for Apache Tomcat -->
<service_bundle type='manifest' name='tomcat9:tomcat9'>
    <service name='network/tomcat9'
             type='service'
             version='1'>
         <create_default_instance enabled='false' />
         <single_instance />
         <dependency name='fs'
                     grouping='require_all'
                     restart_on='none'
                     type='service'>
              <service_fmri value='svc:/system/filesystem/local' />
         </dependency>
         <dependency name='net'
                     grouping='require_all'
                     restart_on='none'
                     type='service'>
              <service_fmri value='svc:/network/loopback' />
           </dependency>
        <!--path /opt/jboss-ews-... will be changed later to reflect real path-->
           <exec_method type='method'
                        name='start'
                        exec='/opt/jboss-ews-2.0/sbin/jws5-tomcat.init start'
                        timeout_seconds='-1'>
           </exec_method>
           <exec_method type='method'
                        name='stop'
                        exec='/opt/jboss-ews-2.0/sbin/jws5-tomcat.init stop'
                        timeout_seconds='-1'>
           </exec_method>
           <exec_method type='method'
                        name='restart'
                        exec='/opt/jboss-ews-2.0/sbin/jws5-tomcat.init restart'
                        timeout_seconds='-1'>
           </exec_method>
    </service>
</service_bundle>
